# roles/vistaicm-server/templates/vistaicm-server.service.j2
[Service]
ExecStart=/usr/bin/node {{install_dir}}/server.js
Restart=always
User={{main_user}}

[Install]
WantedBy=multi-user.target
