# roles/insync/templates/insync.service.j2
[Service]
ExecStart=/opt/insync-portable/client/insync-portable start --no-daemon
Restart=always
User={{main_user}}

[Install]
WantedBy=multi-user.target
